# systemd/prometheus-client-apache-sites.socket
[Socket]
ListenFIFO=%t/prometheus-client-apache-sites.stdin
Service=prometheus-client-apache-sites.service

[Install]
WantedBy=sockets.target